# pihole-monitor/pihole-sync.service
[Unit]
Description=Pi-hole device name sync
After=network-online.target mongod.service pihole-collector.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/rpi-projects/pihole-monitor
ExecStart=/usr/bin/python3 /home/pi/rpi-projects/pihole-monitor/sync.py --watch
Restart=always
RestartSec=15

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
