# pihole-monitor/pihole-collector.service
[Unit]
Description=Pi-hole data collector
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/MonitorPiHole/pihole-monitor
ExecStart=/usr/bin/python3 /home/pi/MonitorPiHole/pihole-monitor/collector.py
Restart=always
RestartSec=10

# Logs disponibles con: journalctl -u pihole-collector -f
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
